# deploy/console.service
[Unit]
Description=Console de pilotage quotidien (Hono + SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=console
Group=console
WorkingDirectory=/opt/console-pilotage
ExecStart=/usr/bin/node server/dist/index.js
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=DB_PATH=/var/lib/console/console.db
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
